# supervisord/apache2.conf
[program:apache2]
command=service apache2 start
numprocs=1
autostart=true
