# supervisord/percona.conf
[program:mysqld]
command=service mysql restart
numprocs=1
autostart=true
